# install/etc/profile.d/golang.sh
export GOPATH=/opt/golibs
export PATH=$PATH:/opt/golang/bin:$GOPATH/bin
